# flask-app.service
[Unit]
Description=Flask app service
After=networking.service

[Service]
WorkingDirectory=/home/wigor/app
ExecStart=/home/wigor/app/venv/